# deploy/legalclear-backend.service
[Unit]
Description=LegalClear FastAPI Backend
After=network.target

[Service]
User=joe
WorkingDirectory=/home/joe/.gemini/antigravity/scratch/legalclear/backend
ExecStart=/home/joe/.gemini/antigravity/scratch/legalclear/backend/venv/bin/python main.py
Restart=always

[Install]
WantedBy=multi-user.target
