# deploy/legalclear-frontend.service
[Unit]
Description=LegalClear React Frontend
After=network.target

[Service]
User=joe
WorkingDirectory=/home/joe/.gemini/antigravity/scratch/legalclear/frontend
ExecStart=/usr/bin/npm run preview -- --port 3000 --host 0.0.0.0
Restart=always

[Install]
WantedBy=multi-user.target
